# tasks/systemd/daily-operations.service
[Unit]
Description=ETF Daily Operations
Documentation=file:///path/to/etf-web/tasks/README.md
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=etf-user
Group=etf-user
WorkingDirectory=/path/to/etf-web
ExecStart=/usr/bin/python3 tasks/daily_operations.py
StandardOutput=append:/path/to/etf-web/logs/daily_operations.log
StandardError=append:/path/to/etf-web/logs/daily_operations_error.log
Environment="PYTHONPATH=/path/to/etf-web"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Restart on failure (optional - only if you want auto-retry)
# Restart=on-failure
# RestartSec=300

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

